# systemd/asistente.service
[Unit]
Description=Asistente Virtual - Orange Pi 5 Ultra
Documentation=https://github.com/tu-usuario/asistente
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/path/to/asistente
Environment="PATH=/path/to/asistente/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/path/to/asistente/venv/bin/python -m src.main
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asistente

# Seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/asistente/logs /path/to/asistente/config

# Limites
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
