# systemd/asistente.service
[Unit]
Description=Asistente Virtual - Orange Pi 5 Ultra
Documentation=https://github.com/tu-usuario/asistente
After=network.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=orangepi
Group=orangepi
WorkingDirectory=/home/orangepi/asistente
Environment="PATH=/home/orangepi/asistente/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/orangepi/asistente/venv/bin/python /home/orangepi/asistente/src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=asistente

# Seguridad
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/orangepi/asistente/logs /home/orangepi/asistente/config

# Limites
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
